# hosts/maelthra/quadlet/ts-ollama.container
# ts-ollama.container — Tailscale sidecar for Ollama
# This container creates the network namespace. ollama joins it via
# Network=ts-ollama.container, so Ollama is directly reachable on the tailnet.
#
# Static IP 10.89.1.3 on app-network lets OpenClaw reach Ollama directly
# at http://10.89.1.3:11434 (permitted by forward_allow in desired-state.yaml).

[Unit]
Description=Tailscale sidecar for Ollama
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Container]
ContainerName=ts-ollama
Image=ghcr.io/tailscale/tailscale:latest
AutoUpdate=registry

Network=app-network.network
IP=10.89.1.3
EnvironmentFile=/etc/homelab-gitops/tailscale-ollama.env

Volume=ts-ollama-state.volume:/var/lib/tailscale:Z
Volume=/etc/homelab-gitops/ts-serve/ollama-serve.json:/config/serve.json:ro,Z

AddDevice=/dev/net/tun:/dev/net/tun
AddCapability=CAP_NET_ADMIN
AddCapability=CAP_NET_RAW

PodmanArgs=--memory=128m --cpus=0.25

[Service]
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target default.target
