# hosts/maelthra/quadlet/ts-openclaw.container
# ts-openclaw.container — Tailscale sidecar for OpenClaw
# Runs at root level for /dev/net/tun access.
# OpenClaw runs as a rootless user-level container and publishes port 18789
# to localhost; Tailscale serve proxies tailnet traffic to that port.

[Unit]
Description=Tailscale sidecar for OpenClaw
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Container]
ContainerName=ts-openclaw
Image=ghcr.io/tailscale/tailscale:latest
AutoUpdate=registry

Network=host
EnvironmentFile=/etc/coreos-gitops/tailscale-openclaw.env

Volume=ts-openclaw-state.volume:/var/lib/tailscale:Z
Volume=/etc/coreos-gitops/ts-serve/openclaw-serve.json:/config/serve.json:ro,Z

# Tailscale needs tun device and NET_ADMIN
AddDevice=/dev/net/tun:/dev/net/tun
AddCapability=CAP_NET_ADMIN
AddCapability=CAP_NET_RAW

NoNewPrivileges=true
ReadOnly=true
Tmpfs=/tmp
Tmpfs=/var/run

PodmanArgs=--memory=128m --cpus=0.25

[Service]
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target default.target
