# hosts/maelthra/systemd/podman-image-prune.service
[Unit]
Description=Prune unused Podman images after auto-update
After=podman-auto-update.service

[Service]
Type=oneshot
ExecStart=podman image prune -f
StandardOutput=journal
StandardError=journal
SyslogIdentifier=podman-image-prune
